# deploy/wa-worker.service
[Unit]
Description=WhatsApp Scheduler - sender worker (neonize)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=raspi
WorkingDirectory=/home/raspi/wa_automation/wa_automation
ExecStart=/home/raspi/wa_automation/venv/bin/python worker.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
